# hadoop_node/lib/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
startsecs=3
autorestart=true
exitcodes=0
redirect_stderr=true
stdout_logfile=/dev/stdout
